# bridge_ip-set.service
[Unit]
Description=Configuração da Interface Bridge
After=network.target
After=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStartPre=/bin/sleep 1
ExecStart=/opt/bridge_ip/bridge_ip-set create_bridge


[Install]
WantedBy=default.target
