# bridge_ip-set.service
[Unit]
Description=Configuração de Interface Bridge
After=network.target
After=NetworkManager.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 15
ExecStart=/opt/bridge_ip/bridge_ip-set

[Install]
WantedBy=default.target

